# fsroot/home/obatiuk/.config/systemd/user/restic-backup-daily@.timer.template
#
# Created by dotfiles setup script on ${TODAY} by ${USER}
#

[Unit]
Description=Daily backup of %i
Requires=network-check.service
After=network-check.service

[Timer]
OnBootSec=15min
RandomizedDelaySec=30min
OnCalendar=daily
Persistent=true
Unit=restic-backup@%i.service

[Install]
WantedBy=timers.target
